# infra/streamer/pricey-mem-watchdog.service
[Unit]
Description=Pricey streamer memory watchdog (pre-emptive recycle before OOM)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pricey-mem-watchdog.sh
