# infra/streamer/pricey-mem-watchdog.timer
[Unit]
Description=Run Pricey streamer memory watchdog every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
